# root/etc/services.d/kerio/run
#!/usr/bin/with-contenv bash
exec s6-setuidgid abc /kerio-connect start
